# edge/gx/service/ovr-dbus2prom/run
#!/bin/sh
set -eu

exec /data/overdrive/dbus2prom/run_exporters.sh
