# edge/systemd/gx-keepalive.timer
[Unit]
Description=Run GX keepalive every 30 seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
